# cimt/com.statcan.cimt-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- GENERATED by make_refresh_plist.py from release_dates.txt. Do not hand-edit.
     CIMT merchandise-trade release days 2026-07 .. 2027-02, 08:35 local (ET). -->
<dict>
  <key>Label</key>
  <string>com.statcan.cimt-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
    <string>/Users/jasonkirby/statcan-explorer/cimt/refresh.py</string>
    <string>--retry</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Month</key><integer>7</integer>
      <key>Day</key><integer>7</integer>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>35</integer>
    </dict>
    <dict>
      <key>Month</key><integer>8</integer>
      <key>Day</key><integer>4</integer>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>35</integer>
    </dict>
    <dict>
      <key>Month</key><integer>9</integer>
      <key>Day</key><integer>3</integer>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>35</integer>
    </dict>
    <dict>
      <key>Month</key><integer>10</integer>
      <key>Day</key><integer>6</integer>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>35</integer>
    </dict>
    <dict>
      <key>Month</key><integer>11</integer>
      <key>Day</key><integer>4</integer>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>35</integer>
    </dict>
    <dict>
      <key>Month</key><integer>12</integer>
      <key>Day</key><integer>8</integer>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>35</integer>
    </dict>
    <dict>
      <key>Month</key><integer>1</integer>
      <key>Day</key><integer>6</integer>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>35</integer>
    </dict>
    <dict>
      <key>Month</key><integer>2</integer>
      <key>Day</key><integer>4</integer>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>35</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/jasonkirby/statcan-explorer/cimt/cimt_refresh.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jasonkirby/statcan-explorer/cimt/cimt_refresh.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
