# gazette-monitor/com.jasonkirby.gazette-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jasonkirby.gazette-digest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/jasonkirby/statcan-explorer/gazette-monitor/send_digest.py</string>
  </array>

  <!-- Fires whenever a digest file lands in the outbox -->
  <key>WatchPaths</key>
  <array>
    <string>/Users/jasonkirby/statcan-explorer/gazette-monitor/outbox</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/jasonkirby/statcan-explorer/gazette-monitor/send_digest.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jasonkirby/statcan-explorer/gazette-monitor/send_digest.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
